# zonely-prod.service
[Unit]
Description=Zonely Production Server with Cloudflare Tunnel
After=network.target docker.service
Wants=network.target
Requires=docker.service

[Service]
Type=simple
User=qingbo
Group=qingbo
WorkingDirectory=/Users/qingbo/Projects/Personal/zonely
Environment=PATH=/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/Users/qingbo
Environment=USER=qingbo

# Load direnv environment and run mix task
ExecStartPre=/bin/bash -c "cd /Users/qingbo/Projects/Personal/zonely && mix db.prod.up"
ExecStart=/bin/bash -c "cd /Users/qingbo/Projects/Personal/zonely && direnv exec . mix prod.tunnel"

# Restart configuration
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Logging
StandardOutput=append:/Users/qingbo/Projects/Personal/zonely/logs/prod.log
StandardError=append:/Users/qingbo/Projects/Personal/zonely/logs/prod.error.log

# Security
NoNewPrivileges=true
PrivateTmp=true

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
